# app/geth-master/templates/geth-master.service.j2
[Unit]
Description=Geth master node service

[Service]
Type=simple
ExecStart=/bin/bash {{ home }}/{{ exec_script }}

[Install]
WantedBy=multi-user.target