# app/geth-slave/templates/geth-slave.service.j2
[Unit]
Description=Geth slave node service

[Service]
Type=simple
ExecStart=/bin/bash {{ home }}/{{ exec_script }}

[Install]
WantedBy=multi-user.target